# root/etc/services.d/sundtek/run
#!/usr/bin/with-contenv bash

cd /opt/bin || exit

exec \
	./mediaclient --start
